# deployment/scanimage.service
[Unit]
Description=ScanImage Web Application
After=network.target

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/scanimage
Environment=PATH=/path/to/your/scanimage/.venv/bin
ExecStart=/path/to/your/scanimage/.venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 4 --daemon --pid /run/scanimage.pid app:app
ExecReload=/bin/kill -s HUP $MAINPID
PIDFile=/run/scanimage.pid
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target